# skills/my-proxy/scripts/com.local.microsocks.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!--
  LaunchAgent for the residential SOCKS5 proxy. Run ONCE on the laptop (macOS).

  Setup:
    1. Replace 100.101.102.103 below with your Tailscale IP:  tailscale ip -4
    2. Copy to ~/Library/LaunchAgents/com.local.microsocks.plist
    3. launchctl load ~/Library/LaunchAgents/com.local.microsocks.plist
    4. launchctl list | grep microsocks   # verify it is running

  Notes:
    - Bound to the Tailscale IP only (microsocks has NO auth) — never use 0.0.0.0.
    - caffeinate keeps the laptop awake while the proxy runs.
    - microsocks path assumes Apple Silicon Homebrew (/opt/homebrew). On Intel use /usr/local.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.local.microsocks</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/caffeinate</string>
    <string>-i</string>
    <string>/opt/homebrew/bin/microsocks</string>
    <string>-i</string>
    <string>100.101.102.103</string>
    <string>-p</string>
    <string>1080</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardErrorPath</key>
  <string>/tmp/microsocks.err.log</string>
  <key>StandardOutPath</key>
  <string>/tmp/microsocks.out.log</string>
</dict>
</plist>
